# src/supervisord.conf
[supervisord]
nodaemon=true

[program:xbmc-server]
command=/opt/xbmc-server/xbmc-server --no-test --nolirc -p
